# docum.service
[Unit]
Description=Docum Service

[Service]
Type=simple
ExecStart=/var/docum/start-server.sh

[Install]
WantedBy=multi-user.target